# airootfs/etc/systemd/system/polyarchy-init.service
[Unit]
Description=Polyarch Bootstrapping Service

[Service]
ExecStart=/etc/systemd/scripts/polyarchy-init

[Install]
RequiredBy=data.automount
